# deploy/openlitespeed/betel-reus.service
[Unit]
Description=Biserica Betel Reus Node.js app
After=network.target

[Service]
Type=simple
User=betel
Group=betel
WorkingDirectory=/var/www/betel-reus
EnvironmentFile=/etc/betel-reus/betel-reus.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
